# docker/gnet.service
[Unit]
Description=Xilinx ml-suite googlenet for deepgreen
After=Netwrok.target

[Service]
ExecStart=/quicksetup/docker/gnet.sh 

